# deploy/systemd/edgelink-call-relayd.service
[Unit]
Description=EdgeLink MiLink-compatible call relay gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/edgelink-call-relayd
ExecStart=/usr/bin/node /opt/edgelink-call-relayd/server.mjs
Restart=always
RestartSec=1
Environment=EDGELINK_CALL_RELAY_PUBLIC_HOST=172.238.24.219
Environment=EDGELINK_CALL_RELAY_CONTROL_PORT=17104
Environment=EDGELINK_CALL_RELAY_PORT_BASE=18000
Environment=EDGELINK_CALL_RELAY_SLOT_COUNT=20
Environment=EDGELINK_WORKER_BASE_URL=https://edgelink-worker.black-hill-f944.workers.dev
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/edgelink-call-relayd

[Install]
WantedBy=multi-user.target
